# systemd/orchestrator.service
[Unit]
Description=Hermes LangGraph Orchestrator
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/
Environment="PYTHONPATH=/orchestrator"
ExecStart=/root/venv/bin/python -m uvicorn orchestrator.orchestrator.api:app --host 127.0.0.1 --port 18793 --log-level info
Restart=on-failure
RestartSec=5
StandardOutput=append:/var/log/orchestrator/stdout.log
StandardError=append:/var/log/orchestrator/stderr.log

[Install]
WantedBy=multi-user.target
