# systemd/orchestrator.service
[Unit]
Description=LangGraph Orchestrator Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/orchestrator
Environment="PYTHONPATH=/orchestrator"
EnvironmentFile=/config/orchestrator.env
ExecStart=/root/venv/bin/python -m uvicorn orchestrator.api:app --host 127.0.0.1 --port 18793 --log-level info
Restart=on-failure
RestartSec=5
StandardOutput=append:/var/log/orchestrator/stdout.log
StandardError=append:/var/log/orchestrator/stderr.log

[Install]
WantedBy=multi-user.target
